# results/words-xaq.list
binderhubovh/ovhbhub-fumitoh-2dmls2018-2ddemo-c32019; [['-z', '$(apt-get indextargets)'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'ca-certificates', 'curl', 'netbase', 'wget'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'bzr', 'git', 'mercurial', 'openssh-client', 'subversion', 'procps'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', '-qq', 'update'], ['apt-get', '-qq', 'purge'], ['apt-get', '-qq', 'clean'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['locale-gen']]
futuretechindustries/pgbouncer-kube; [['rm', '-rf', '/var/lib/apt/lists/*'], ['sed', '-i', 's/^#s*(deb.*universe)$/1/g', '/etc/apt/sources.list'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt', 'install', '-y', 'pgbouncer'], ['adduser', 'pgbouncer', '--no-create-home']]
awsservicebroker/dynamodb-apb; [['curl', 'https://copr.fedorainfracloud.org/coprs/g/ansible-service-broker/ansible-service-broker-latest/repo/epel-7/group_ansible-service-broker-ansible-service-broker-latest-epel-7.repo', '-o', '/etc/yum.repos.d/asb.repo'], ['yum', '-y', 'install', 'epel-release', 'centos-release-openshift-origin'], ['yum', '-y', 'install', '--setopt=tsflags=nodocs', 'origin-clients', 'python-openshift', 'ansible', 'ansible-kubernetes-modules', 'ansible-asb-modules', 'apb-base-scripts'], ['yum', 'clean', 'all'], ['mkdir', '-p', '/usr/share/ansible/openshift', '/etc/ansible', '/opt/ansible', '${BASE_DIR}/{etc,.kube,.ansible/tmp}'], ['useradd', '-u', '${USER_UID}', '-r', '-g', '0', '-M', '-d', '${BASE_DIR}', '-b', '${BASE_DIR}', '-s', '/sbin/nologin', '-c', 'apb user', '${USER_NAME}'], ['chown', '-R', '${USER_NAME}:0', '/opt/{ansible,apb}'], ['chmod', '-R', 'g=u', '/opt/{ansible,apb}', '/etc/passwd'], ['yum', '-y', 'install', 'python-boto3', 'python-netaddr'], ['yum', 'clean', 'all'], ['chmod', '+x', '/usr/bin/16to62'], ['chmod', '+x', '/usr/bin/ssmfetch']]
jisto/ubuntu1404-puppet; [['rm', '-rf', '/var/lib/apt/lists/*'], ['sed', '-i', 's/^#s*(deb.*universe)$/1/g', '/etc/apt/sources.list'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt-get', 'install', '-y', 'software-properties-common', 'python-software-properties', 'time', 'vim', 'git', 'tmux', 'curl', 'systemd', 'iptables'], ['cd', '/lib/systemd/system/sysinit.target.wants/']]
jisto/centos7-puppet; [['yum', 'install', '-y', 'tar', 'net-tools', 'git', 'vim', 'strace', 'lsof', 'sudo', 'curl'], ['cd', '/lib/systemd/system/sysinit.target.wants/']]
binderhubovh/ovhbhub-https-3a-2f-2fbitbucket-2eorg-2fbrezyl-2fcofr-2fsrc-2fmaster-c22f42; [['-z', '$(apt-get indextargets)'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'ca-certificates', 'curl', 'netbase', 'wget'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'bzr', 'git', 'mercurial', 'openssh-client', 'subversion', 'procps'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', '-qq', 'update'], ['apt-get', '-qq', 'purge'], ['apt-get', '-qq', 'clean'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['locale-gen']]
binderhubovh/ovhbhub-https-3a-2f-2frisk-2dengineering-2eorg-2fgit-2fnotebooks-2egit-6829f8; [['-z', '$(apt-get indextargets)'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'ca-certificates', 'curl', 'netbase', 'wget'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'bzr', 'git', 'mercurial', 'openssh-client', 'subversion', 'procps'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', '-qq', 'update'], ['apt-get', '-qq', 'purge'], ['apt-get', '-qq', 'clean'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['locale-gen']]
binderhubovh/ovhbhub-https-3a-2f-2fframagit-2eorg-2fdebimax-2fcours-2ddebimax-a31fab; [['-z', '$(apt-get indextargets)'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'ca-certificates', 'curl', 'netbase', 'wget'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'bzr', 'git', 'mercurial', 'openssh-client', 'subversion', 'procps'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', '-qq', 'update'], ['apt-get', '-qq', 'purge'], ['apt-get', '-qq', 'clean'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['locale-gen']]
binderhubovh/ovhbhub-hunters-2dforge-2dmordor-cf2bad; [['rm', '-rf', '/var/lib/apt/lists/*'], ['sed', '-i', 's/^#s*(deb.*universe)$/1/g', '/etc/apt/sources.list'], ['mkdir', '-p', '/run/systemd'], ['/bd_build/prepare.sh'], ['/bd_build/system_services.sh'], ['/bd_build/utilities.sh'], ['/bd_build/cleanup.sh']]
binderhubovh/ovhbhub-hackingmaterials-2dmatminer-6f2b9d; [['-z', '$(apt-get indextargets)'], ['mkdir', '-p', '/run/systemd'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'ca-certificates', 'curl', 'netbase', 'wget'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', 'update'], ['apt-get', 'install', '-y', '--no-install-recommends', 'bzr', 'git', 'mercurial', 'openssh-client', 'subversion', 'procps'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['apt-get', '-qq', 'update'], ['apt-get', '-qq', 'purge'], ['apt-get', '-qq', 'clean'], ['rm', '-rf', '/var/lib/apt/lists/*'], ['locale-gen']]
